# api/api-supervisord.conf
[program:uwsgi-api]
command = /usr/local/bin/uwsgi --ini /opt/uwsgi.ini

[program:nginx-api]
command = /usr/sbin/nginx